# systemd/tubpi-webapp.service
[Unit]
Description=Tubpi Web App - Interface web de contrôle du rail
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=gpio
WorkingDirectory=/opt/tubpi/src
Environment="PYTHONUNBUFFERED=1"
Environment="FLASK_APP=web_app.py"
ExecStart=/usr/bin/python3 /opt/tubpi/src/web_app.py
Restart=always
RestartSec=10

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tubpi-webapp

# Sécurité
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
