# systemd/tubpi-webapp.service
[Unit]
Description=Tubpi Web App - Interface web de contrôle du rail
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=gpio
SupplementaryGroups=gpio
WorkingDirectory=/opt/tubpi/src
Environment="PYTHONUNBUFFERED=1"
Environment="FLASK_APP=web_app.py"
ExecStart=/usr/bin/python3 /opt/tubpi/src/web_app.py
Restart=always
RestartSec=10

# Permissions GPIO pour Raspberry Pi 5
DeviceAllow=/dev/gpiochip0 rw
DeviceAllow=/dev/gpiochip1 rw
DeviceAllow=/dev/gpiochip2 rw
DeviceAllow=/dev/gpiochip3 rw
DeviceAllow=/dev/gpiochip4 rw
DeviceAllow=/dev/gpiomem rw
DevicePolicy=closed

# Logs
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tubpi-webapp

# Sécurité
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
